# deploy/politest.service
# systemd unit for the DecaCompass crowd-comparison API.
# Install: sudo cp /opt/politest/deploy/politest.service /etc/systemd/system/
#          sudo systemctl daemon-reload && sudo systemctl enable --now politest
#          journalctl -u politest -f      # confirm it started
[Unit]
Description=Politeion API (crowd comparison + article analyzer)
After=network.target

[Service]
Type=simple
User=profileher
Group=profileher
WorkingDirectory=/opt/politest
ExecStart=/usr/bin/node server.js
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=3200
Environment=STORE_FILE=/opt/politest/store/results.jsonl
Environment=ANALYSES_FILE=/opt/politest/store/analyses.jsonl
Environment=USAGE_FILE=/opt/politest/store/analyzer-usage.jsonl
# Analyzer secrets/config (ANTHROPIC_API_KEY, PROVIDER, MODEL, MONTHLY_BUDGET_USD).
# systemd (root) reads this at start and injects the vars, then drops to User=
# above — so the file stays root:root 0600 and the app process never reads it.
# The leading "-" keeps the service starting if the file is absent (the Analyzer
# then reports "not configured" and refuses jobs, while the rest of the app runs).
EnvironmentFile=-/etc/politeion/analyzer.env
Restart=on-failure
RestartSec=3

# Hardening. ProtectSystem=strict makes the FS read-only except ReadWritePaths,
# which is where the append-only crowd store lives.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/politest/store

[Install]
WantedBy=multi-user.target
