# deploy/politest-brief.service
# Oneshot: draft (never publish) the Politeion Daily Brief. Runs as the same
# unprivileged user as the API, reads the same EnvironmentFile (ANTHROPIC_API_KEY,
# MODEL, MONTHLY_BUDGET_USD). Drafting NEVER auto-publishes — an admin approves each
# brief from account.html. If the month's budget cap is hit, the draft aborts cleanly
# (exit 0) and produces nothing.
#
# Install:
#   sudo cp /opt/politest/deploy/politest-brief.service /etc/systemd/system/
#   sudo cp /opt/politest/deploy/politest-brief.timer   /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now politest-brief.timer
#   systemctl list-timers politest-brief.timer          # confirm next run
#   sudo systemctl start politest-brief.service         # draft now (manual test)
[Unit]
Description=Politeion Daily Brief — draft (manual approval required to publish)
After=network.target

[Service]
Type=oneshot
User=profileher
Group=profileher
WorkingDirectory=/opt/politest
ExecStart=/usr/bin/node tools/brief-draft.js
Environment=NODE_ENV=production
Environment=ANALYSES_FILE=/opt/politest/store/analyses.jsonl
Environment=USAGE_FILE=/opt/politest/store/analyzer-usage.jsonl
Environment=BRIEFS_FILE=/opt/politest/store/briefs.jsonl
EnvironmentFile=-/etc/politeion/analyzer.env

# Same hardening as the API: read-only FS except the writable store.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/politest/store
